# deploy/telegrambot.service
[Unit]
Description=Telegram Broadcast Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
EnvironmentFile=__WORKDIR__/.env
ExecStart=__PYTHON__ __WORKDIR__/main.py
User=__USER__
Group=__USER__
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
